# systemd/photoframe.service
# Systemd Service Template for Photo Frame (Headless DRM/KMS)
# Optimized for Raspberry Pi OS Trixie headless operation
# Placeholders are replaced during installation by install.sh

[Unit]
Description=Pi3D Photo Frame Headless (DRM/KMS)
Documentation=https://github.com/frankCls/photo-frame
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi                      # TEMPLATE
Group=pi                     # TEMPLATE
WorkingDirectory=/home/pi/photo-frame/processed_photos  # TEMPLATE

# Environment variables
Environment="HOME=/home/pi"  # TEMPLATE
Environment="PYTHONPATH=/home/pi/pi3d_demos"  # TEMPLATE
Environment="SDL_VIDEODRIVER=kmsdrm"
Environment="LANG=en_US.UTF-8"
Environment="LC_ALL=en_US.UTF-8"

# Direct Python execution (no compositor)
ExecStart=__EXEC_START__  # TEMPLATE

# GPU/KMS access permissions
SupplementaryGroups=video render input tty

# Restart on failure
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=photoframe

[Install]
WantedBy=multi-user.target
