# systemd/photoframe.service
# Systemd Service Template for Photo Frame (cage compositor)
# Optimized for Raspberry Pi OS Trixie with Wayland/KMS
# Placeholders are replaced during installation by install.sh

[Unit]
Description=Raspberry Pi Photo Frame (Pi3D Slideshow)
Documentation=https://github.com/frankCls/photo-frame
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi                      # TEMPLATE
Group=pi                     # TEMPLATE
WorkingDirectory=/home/pi/photo-frame/processed_photos  # TEMPLATE

# Environment variables
Environment="HOME=/home/pi"  # TEMPLATE
Environment="XDG_RUNTIME_DIR=/run/user/__UID__"  # TEMPLATE

# TTY and PAM settings for hardware access
PAMName=login
TTYPath=/dev/tty1
StandardInput=tty

# Supplementary groups for GPU/KMS access
SupplementaryGroups=video render input tty

# Create XDG_RUNTIME_DIR if needed
RuntimeDirectory=user/__UID__  # TEMPLATE

# cage command: /usr/bin/cage -s -- python3 PictureFrame2020.py -p /photos
ExecStart=__EXEC_START__  # TEMPLATE

# Restart on failure
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=photoframe

[Install]
WantedBy=graphical.target
